# scripts/armor-version-drift-check.cron
# ARMOR Version Drift Check - Cron Job
# Check daily at 9:17 AM (avoiding :00 and :30 marks)
17 9 * * * /home/coding/ARMOR/scripts/check-armor-version-drift.py > /home/coding/ARMOR/logs/version-drift-check.log 2>&1
